# configs/bms.timer
[Unit]
Description=Run bms service every 20 min

[Timer]
OnCalendar=*:0/20
